# com.investment-agent.morning-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.investment-agent.morning-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/juntanaka/Applications/investment-info/venv/bin/python</string>
    <string>/Users/juntanaka/Applications/investment-info/main.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ANTHROPIC_API_KEY</key>
    <string>YOUR_ANTHROPIC_API_KEY</string>
    <key>SLACK_WEBHOOK_URL</key>
    <string>YOUR_SLACK_WEBHOOK_URL</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/juntanaka/Applications/investment-info/logs/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/juntanaka/Applications/investment-info/logs/stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
